# debian/server/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Reload firewalld if it's running
        if systemctl is-active --quiet firewalld 2>/dev/null; then
            firewall-cmd --reload >/dev/null 2>&1 || true
        fi

        echo "PulseAudio Network Server installed successfully!"
        echo ""
        echo "To use this service:"
        echo "1. Enable and start the user service:"
        echo "   systemctl --user enable pulseaudio-network-server.service"
        echo "   systemctl --user start pulseaudio-network-server.service"
        echo ""
        echo "2. Configure firewall (as root):"
        echo "   firewall-cmd --permanent --add-service=pulseaudio-network"
        echo "   firewall-cmd --reload"
        echo ""
        echo "3. Edit configuration if needed:"
        echo "   ~/.config/pulseaudio-network/server.json"
        ;;
esac

#DEBHELPER#

exit 0