# debian/client/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        echo "PulseAudio Network Client installed successfully!"
        echo ""
        echo "Configuration Management:"
        echo "1. Create/edit configurations:"
        echo "   pulseaudio-network-client-config create-sample"
        echo "   pulseaudio-network-client-config list"
        echo ""
        echo "2. Enable specific configurations:"
        echo "   pulseaudio-network-client-config enable default"
        echo "   pulseaudio-network-client-config enable office"
        echo ""
        echo "3. Using systemd directly:"
        echo "   systemctl --user enable pulseaudio-network-client@office.service"
        echo "   systemctl --user start pulseaudio-network-client@office.service"
        echo ""
        echo "4. Check service status:"
        echo "   pulseaudio-network-client-config status office"
        ;;
esac

#DEBHELPER#

exit 0
